# vswm.service
[Unit]
After=networking.target

[Service]
Type=forking
ExecStart=/usr/bin/python3 /usr/local/bin/vswm

[Install]
WantedBy=networking.target
